# config/profile.d/ioc.sh
export IOC_CONFIG_FILE=/etc/ioc/config.toml
